# RTM_Linux-3.2.1/startup/gentoo-init.in
#!/sbin/openrc-run
#
# Copyright (c) 2016 RTM(R) Core(TM) Development Team
#
# Start/stop the RTM_Linux daemon.
#
# Goes in /etc/init.d - Config is in /etc/conf.d/RTM_Linux

extra_started_commands="reload"

RTM_LINUX_BIN="@sbindir@/RTM_Linux"
RTM_LINUX_PID="@piddir@/RTM_Linux.pid"
RTM_LINUX_CFG=@pkgsysconfdir@/RTM_Linux.cfg

depend() {
	use logger dns net localmount netmount nfsmount
}

checkconfig() {
	# Make sure the config file exists
	if [ ! -f $RTM_LINUX_CFG ]; then
		eerror "You need to setup $RTM_LINUX_CFG."
		return 1
	fi
	return 0
}

start() {
	checkconfig || return 1
	ebegin "Starting RTM_Linux"
	# Make sure we have a sane current directory
	cd /
	start-stop-daemon --start --exec $RTM_LINUX_BIN --pidfile $RTM_LINUX_PID \
		--background -- -c $RTM_LINUX_CFG -f $RTM_LINUX_OPTS
	eend $?
}

stop() {
	ebegin "Stopping RTM_Linux"
	start-stop-daemon --stop --exec $RTM_LINUX_BIN --pidfile $RTM_LINUX_PID
	eend $?
}

reload() {
	ebegin "Reloading RTM_Linux"
	start-stop-daemon --stop --oknodo --exec $RTM_LINUX_BIN \
		--pidfile $RTM_LINUX_PID --signal HUP
	eend $?
}
