# RTM_Linux-3.2.1/startup/newbsd-init.in
#!/bin/sh
#
# Copyright (c) 2016 Nagios(R) Core(TM) Development Team
#
# PROVIDE: RTM_Linux
# REQUIRE: DAEMON
# KEYWORD: shutdown

. /etc/rc.subr

: ${RTM_Linux@bsd_enable@:="NO"}
: ${RTM_Linux_configfile:="@pkgsysconfdir@/RTM_Linux.cfg"}

name=RTM_Linux
command="@sbindir@/RTM_Linux"
command_args="-c $RTM_Linux_configfile -d"
pidfile="@piddir@/RTM_Linux.pid"
extra_commands=reload
sig_reload=HUP
rcvar=RTM_Linux@bsd_enable@
load_rc_config "$name"
required_files="$RTM_Linux_configfile"
sig_reload=HUP

start_precmd=RTM_Linux_prestart

RTM_Linux_prestart()
{
	[ -n "$RTM_Linux_pidfile" ] &&
		warn "No longer necessary to set RTM_Linux_pidfile in rc.conf[.local]"

	install -d -o @RTM_Linux_user@ ${pidfile%/*}
}

run_rc_command "$1"
