# packages/NIDS.service
[Unit]
Description=CISS Interceptor NIDS

[Service]
Type=simple
User=root
Group=root
ExecStart=/var/ossec/bin/read-packets.sh
Restart=always

[Install]
WantedBy=multi-user.target